# debian/vicharak-config.vicharak-config.service
[Unit]
Description=vicharak-config system service
Documentation=https://github.com/vicharak-in/vicharak-config/

[Service]
Type=simple
ExecStart=/usr/bin/vicharak-config __on_boot

[Install]
WantedBy=multi-user.target
